# services/ledplayer-kiosk.service
[Unit]
Description=LedscreenPlayer Kiosk Browser (LED region window)
After=graphical.target ledplayer-api.service
Wants=ledplayer-api.service

[Service]
Type=simple
User=pi
Environment=DISPLAY=:0
Environment=XAUTHORITY=/home/pi/.Xauthority

# Wait for the API server and X session to be ready
ExecStartPre=/bin/sleep 5

ExecStart=/bin/bash -c '\
  CHROMIUM=$(command -v chromium || command -v chromium-browser) && \
  $CHROMIUM \
    --kiosk \
    --app=http://localhost:8000/player \
    --noerrdialogs \
    --disable-infobars \
    --no-first-run \
    --disable-translate \
    --disable-features=TranslateUI \
    --disable-session-crashed-bubble \
    --disable-restore-session-state \
    --autoplay-policy=no-user-gesture-required \
    --check-for-update-interval=31536000 \
    --disable-background-timer-throttling \
    --disable-backgrounding-occluded-windows \
    --disable-renderer-backgrounding \
    --hide-scrollbars \
    --overscroll-history-navigation=0'

Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ledplayer-kiosk

[Install]
WantedBy=graphical.target
